# additions/mounts_issaf/data.automount
[Unit]
Description=Mount NFS /data

[Automount]
Where=/data
DirectoryMode=0777
##TimeoutSec=0

[Install]
WantedBy=multi-user.target

