# additions/mounts_issaf/data6.automount
[Unit]
Description=Mount NFS data6

[Automount]
Where=/data6
DirectoryMode=0777
##TimeoutSec=0

[Install]
WantedBy=multi-user.target

