# additions/mounts_issaf/fastdata.automount
[Unit]
Description=Mount NFS /fastdata

[Automount]
Where=/fastdata
DirectoryMode=0777
##TimeoutSec=0

[Install]
WantedBy=multi-user.target

