# additions/mounts_issaf/share-apps.automount
[Unit]
Description=Mount NFS /share/apps

[Automount]
Where=/share/apps
DirectoryMode=0777
##TimeoutSec=0

[Install]
WantedBy=multi-user.target

