# additions/mounts_alien/export.automount
[Unit]
Description=AutoMount NFS /export

[Automount]
Where=/export
DirectoryMode=0755
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target

