# additions/mounts_vmhosts/vmhost2/data02.automount
[Unit]
Description=AutoMount NFS /data02

[Automount]
Where=/data02
DirectoryMode=0755
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target

