# additions/eos_setup/eos@.service.d/override.conf
[Service]
TimeoutStartSec=300

[Service]
ExecStart=
ExecStart=/opt/eos/xrootd/bin/xrootd -n %i -c /etc/xrd.cf.%i -l /var/log/eos/xrdlog.%i -Rdaemon

