# additions/mounts_alien/ini/var-arc-sessions.automount
[Unit]
Description=Mount NFS home
ConditionHost=!alien.spacescience.ro

[Automount]
Where=/var/arc/sessions
DirectoryMode=0777
TimeoutSec=0

[Install]
WantedBy=multi-user.target

