# additions/services_vmhosts/esp_sync.timer
[Unit]
Description=Run esp_sync

[Timer]
OnStartupSec=30

[Install]
WantedBy=timers.target

